# crontab
#run python script every day at midnight
0 0 * * * /usr/local/bin/python /github/app/main.py >> /github/logs/github_backup.log 2>&1
